# scripts/ipc-to-be-arranged/gpslog/gpslog.service
[Unit]
Description=GPSlog

[Service]
Type=Simple
ExecStart=/bin/bash /opt/gpslog.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
