# scripts/ipc-to-be-arranged/mount_files/mnt-disk_C1.automount
Description="Automount External Drive"

[Automount]
Where=/mnt/disk_C1

[Install]
WantedBy=multi-user.target

